# docker/php/supervisord.conf
; Supervisord manages multiple processes in a single container.
;
; Why not just run php-fpm?
; Because we also need a queue worker running alongside it.
; The queue worker picks up async jobs from RabbitMQ (search sync, price alerts).
; Supervisord keeps both processes alive and restarts them if they crash.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; Process 1: PHP-FPM
; This handles incoming HTTP requests from Nginx (GraphQL queries).
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Process 2: Laravel Queue Worker
; This listens to RabbitMQ and processes jobs asynchronously.
; --queue=search-sync,notifications,recommendations defines priority order:
;   1. search-sync: keeping Elasticsearch up to date (highest priority)
;   2. notifications: sending price alerts to users
;   3. recommendations: calculating similar phones (lowest priority, can wait)
; --sleep=3: wait 3 seconds between polling when the queue is empty
; --tries=3: retry a failed job up to 3 times before marking it as failed
[program:queue-worker]
command=php /var/www/html/artisan queue:work rabbitmq --queue=search-sync,notifications,recommendations --sleep=3 --tries=3
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
